# crontab
* * * * * python3 ~/routines/import/metno/hourly/model.py
